# scripts/pingtest.conf
description     "Ping test"
author          "Jason Sayre"

start on startup
stop on runlevel [016]

respawn


exec /vagrant/pingtest.sh
